# provisioning/roles/common/templates/auth-server.service.j2
[Unit]
Description=Authorization Service

[Service]
User={{username}}
TimeoutStartSec=infinity
Environment=ISSUER_URL={{auth_server_base}}/
Environment=FHIR_SERVER_URL={{fhir_server_base}}
Environment=APP_SERVER_URL={{app_server_base}}
Environment=MYSQL_DB={{auth_server_database}}
Environment=MYSQL_USER={{mysql_user}}
Environment=MYSQL_PASS={{mysql_pass}}
Environment=LOGOUT_SUCCESS_URL_PROPERTY={{auth_server_base}}
WorkingDirectory={{install_dir}}/hspc/reference-auth-server-ldap-webapp/target
{% if auth_server_secure_http and not use_custom_ssl_certificates %}
ExecStart=/usr/bin/java -Djetty.port="9001" -Djetty.host="localhost" -Djavax.net.ssl.trustStore="{{install_dir}}/keystore" -Djavax.net.ssl.trustStorePassword="{{keystore_password}}" -jar {{install_dir}}/jetty-runner.jar --config {{install_dir}}/jetty.xml hspc-reference-authorization.war
{% else %}
ExecStart=/usr/bin/java -Djetty.port="9001" -Djetty.host="localhost" -jar {{install_dir}}/jetty-runner.jar --config {{install_dir}}/jetty.xml hspc-reference-authorization.war
{% endif %}
Restart=always

[Install]
WantedBy=multi-user.target
